# meta-balena-common/recipes-connectivity/dnsmasq/balena-files/dnsmasq-servers.path
[Unit]
Description=dnsmasq.servers path watch
Requires=resin-boot.service
After=resin-boot.service

[Path]
PathChanged=/run/dnsmasq.servers

[Install]
WantedBy=basic.target
